# iot/rc.local
#!/bin/bash

/etc/iot/startup.sh >> "/etc/iot/startup.log" 2>&1

exit 0
